# aiwholesail-api/scripts/systemd/aiwholesail-onfailure-notify.service
[Unit]
Description=Email operator when an aiwholesail-* unit fails
# Don't propagate failures from this notifier — would create alert loop
# if it itself crashed.

[Service]
Type=oneshot
User=root
WorkingDirectory=/var/www/aiwholesail-api
EnvironmentFile=/var/www/aiwholesail-api/.env
TimeoutStartSec=60
# MONITOR_UNIT, MONITOR_INVOCATION_ID, MONITOR_EXIT_CODE are set by systemd
# automatically when this unit is invoked via OnFailure= from another unit.
ExecStart=/usr/bin/node scripts/onfailure-notify.js
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aiwholesail-onfailure-notify
Nice=10

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
LockPersonality=true
RestrictRealtime=true
SystemCallArchitectures=native
ReadWritePaths=/var/www/aiwholesail-api /var/log
