# aiwholesail-api/scripts/systemd/aiwholesail-rapidapi-alert.service
[Unit]
Description=AIWholesail RapidAPI alert worker (evaluates SLO 9/10 + emits alerts)
After=network-online.target postgresql.service aiwholesail-api.service
Wants=network-online.target postgresql.service
OnFailure=aiwholesail-onfailure-notify.service

[Service]
Type=oneshot
User=root
WorkingDirectory=/var/www/aiwholesail-api
EnvironmentFile=/var/www/aiwholesail-api/.env
ExecStart=/usr/bin/node scripts/rapidapi-alert-worker.js
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aiwholesail-rapidapi-alert
Nice=10

# Hardening (matches main API unit)
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
LockPersonality=true
RestrictRealtime=true
SystemCallArchitectures=native
ReadWritePaths=/var/www/aiwholesail-api /var/log
