# aiwholesail-api/scripts/systemd/aiwholesail-founder-hot-list.service
[Unit]
Description=AIWholesail founder hot-list worker (curated daily trial outreach list for Connor)
After=network-online.target postgresql.service aiwholesail-api.service
Wants=network-online.target postgresql.service
OnFailure=aiwholesail-onfailure-notify.service

[Service]
Type=oneshot
User=root
WorkingDirectory=/var/www/aiwholesail-api
EnvironmentFile=/var/www/aiwholesail-api/.env
ExecStart=/usr/bin/node scripts/founder-hot-list.js
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aiwholesail-founder-hot-list
Nice=10

# Hardening (matches main API unit)
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
LockPersonality=true
RestrictRealtime=true
SystemCallArchitectures=native
ReadWritePaths=/var/www/aiwholesail-api /var/log
