# aiwholesail-api/scripts/systemd/aiwholesail-offmarket-routing-monitor.service
[Unit]
Description=AIWholesail off-market routing-collapse monitor (SLI checks + Resend alerts)
After=network-online.target postgresql.service
Wants=network-online.target postgresql.service

[Service]
Type=oneshot
# journalctl access requires either the `systemd-journal` group or root.
# Mirroring health-monitor.service which also reads journald.
User=root
WorkingDirectory=/var/www/aiwholesail-api
EnvironmentFile=/var/www/aiwholesail-api/.env
TimeoutStartSec=120
ExecStart=/usr/bin/node scripts/offmarket-routing-monitor.js
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aiwholesail-offmarket-monitor
Nice=15

# Hardening — matches other aiwholesail-* one-shots
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
LockPersonality=true
RestrictRealtime=true
SystemCallArchitectures=native
ReadWritePaths=/var/www/aiwholesail-api /var/log
