# aiwholesail-api/scripts/systemd/aiwholesail-trial-expiry-sweep.service
[Unit]
Description=AIWholesail trial-expiry sweep (downgrades abandoned trials past trial_end)
After=network-online.target postgresql.service aiwholesail-api.service
Wants=network-online.target postgresql.service
OnFailure=aiwholesail-onfailure-notify.service

[Service]
Type=oneshot
User=root
WorkingDirectory=/var/www/aiwholesail-api
EnvironmentFile=/var/www/aiwholesail-api/.env
ExecStart=/usr/bin/node scripts/trial-expiry-sweep.js
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aiwholesail-trial-expiry-sweep
Nice=10

# Hardening (matches trial-lifecycle unit)
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
LockPersonality=true
RestrictRealtime=true
SystemCallArchitectures=native
ReadWritePaths=/var/www/aiwholesail-api /var/log
